# automation/com.stock.before_market.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- 盤前分析自動化任務 -->
    <key>Label</key>
    <string>com.stock.before_market</string>

    <!-- 執行的程式 -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/walter/Documents/GitHub/stock/automation/run_before_market.py</string>
    </array>

    <!-- 工作目錄 -->
    <key>WorkingDirectory</key>
    <string>/Users/walter/Documents/GitHub/stock</string>

    <!-- 執行時間：每天早上08:00 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- 標準輸出日誌 -->
    <key>StandardOutPath</key>
    <string>/Users/walter/Documents/GitHub/stock/logs/before_market.log</string>

    <!-- 錯誤輸出日誌 -->
    <key>StandardErrorPath</key>
    <string>/Users/walter/Documents/GitHub/stock/logs/before_market.error.log</string>

    <!-- 失敗後自動重啟 -->
    <key>KeepAlive</key>
    <false/>

    <!-- 開機自動啟動 -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
